# whatsapp-automation.service
[Unit]
Description=WhatsApp Automation Docker Containers
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ubuntu/wa-auto-v2
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=0
User=ubuntu
Group=docker

[Install]
WantedBy=multi-user.target